# deploy/systemd-user/english-tech-postgres.service
[Unit]
Description=English Tech local PostgreSQL
After=network.target

[Service]
Type=forking
ExecStart=/home/aseda/Desktop/english_tech/scripts/local_services/start_local_postgres.sh
ExecStop=/home/aseda/Desktop/english_tech/scripts/local_services/stop_local_postgres.sh
Restart=on-failure

[Install]
WantedBy=default.target
